# rclone/.config/systemd/user/rclone-gdrive.service
[Unit]
Description=Sync files with GDrive

[Service]
Type=oneshot
ExecStart=rclone bisync %h/GDrive GDrive:/

[Install]
WantedBy=default.target
