# rclone/.config/systemd/user/rclone-gdrive.timer
[Unit]
Description=Ensure GDrive is in sync

[Timer]
OnBootSec=1m
OnUnitActiveSec=10m

[Install]
WantedBy=timers.target
