# rclone/.config/systemd/user/rclone-gdrive.timer
[Unit]
Description=Ensure GDrive is in sync

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
