# rclone/.config/systemd/user/rclone-gdrive.path
[Unit]
Description=Watch GDrive files for changes and trigger rclone sync

[Path]
PathChanged=%h/GDrive/
Unit=rclone-gdrive.service

[Install]
WantedBy=default.target
