# emacs/.config/systemd/user/emacs.service.d/override.conf
[Service]
ExecStop=
ExecStop=/usr/bin/emacsclient --eval "(server-save-buffers-kill-terminal t)"
Type=simple
