# emacs/.config/systemd/user/emacs.service.d/override.conf
[Service]
ExecStop=
ExecStop=/usr/bin/emacsclient --eval "(let ((current-prefix-arg '(4))) (call-interactively 'save-buffers-kill-emacs))"
