# wallplug/.config/systemd/user/wallplug.service
[Unit]
Description=Periodically set wallpaper

[Service]
ExecStart=%h/.local/bin/wallplug.sh