# wallplug/.config/systemd/user/wallplug.service
[Unit]
Description=Periodically set wallpaper
PartOf=graphical-session.target

[Service]
ExecStart=%h/.local/bin/wallplug.sh
