# disk-monitor/.config/systemd/user/disk-monitor.timer
[Unit]
Description=Notify low disk space

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
