# sudo-notify/.config/systemd/user/sudo-notify.timer
[Unit]
Description=Notify sudo request (eventually)

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
