# sudo-notify/.config/systemd/user/sudo-notify.service
[Unit]
Description=Notify sudo request (eventually)

[Service]
ExecStartPre=-killall sudo-notify.py
ExecStart=%h/.local/bin/sudo-notify.py
